# packaging/autostart/gpu-monitor-tray.desktop
[Desktop Entry]
Type=Application
Name=GPU Monitor Tray
Comment=System tray indicator for GPU usage (consumes gpu-monitord)
Exec=sh -c "exec $HOME/.local/bin/gpu-monitor-tray"
Icon=tarjeta-de-video
Terminal=false
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Categories=System;Monitor;
